# AXE-CORE-ORCHESTRATOR-content/AXE-CORE-HEADQUARTERS/backend/axe_api/axe_api.service
[Unit]
Description=AXE Core API
Documentation=https://github.com/Ldezeeuw445/AXE-CORE-
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/axe-core-api
EnvironmentFile=/opt/axe-core-api/.env
ExecStart=/opt/axe-core-api/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
